# webtranscriber-main.service
[Unit]
Description=Web Transcriber API Server (main.py)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/web-transcriber
ExecStart=/opt/web-transcriber/.venv/bin/python main.py --backend
Restart=on-failure
User=webtranscriber

[Install]
WantedBy=multi-user.target
