# webtranscriber-webui.service
[Unit]
Description=Web Transcriber Web UI (webui.py)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/rxon/projects/web-transcriber
Environment="PATH=/home/rxon/projects/web-transcriber/.venv/bin:/usr/bin:/bin"
ExecStart=/home/rxon/projects/web-transcriber/.venv/bin/uvicorn webui:app --host 0.0.0.0 --port 8001 --reload
Restart=on-failure
User=rxon

[Install]
WantedBy=multi-user.target
