# webtranscriber-webui.service
[Unit]
Description=Web Transcriber Web UI (webui.py)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/web-transcriber
ExecStart=/opt/web-transcriber/.venv/bin/python main.py --frontend
Restart=on-failure
User=webtranscriber

[Install]
WantedBy=multi-user.target
